# scripts/event-tracker-backup.service
[Unit]
Description=Personal Event Tracker — SQLite backup
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/docker/pet
ExecStart=/opt/docker/pet/scripts/backup-db.sh
Nice=10
IOSchedulingClass=idle
